# cj-gwellcity1-addon/rootfs/etc/services.d/gwellcity1/run
#!/usr/bin/with-contenv bashio
bashio::log.info "Starting CJ Gwellcity1 service"
exec python3 /usr/bin/gwellcity1.py
